# systemd/sentinel.service
[Unit]
Description=Sentinel Defense Stock Claim Analyzer API
After=network-online.target postgresql.service
Wants=network-online.target

StartLimitIntervalSec=1800
StartLimitBurst=5

[Service]
Type=exec
User=sentinel
Group=sentinel
WorkingDirectory=/opt/sentinel/app

ExecStart=/home/sentinel/.local/bin/uv run serve --port 5000 --workers 4

Restart=on-failure
RestartSec=30

# Kill the entire process group (master + workers), not just the master
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

MemoryMax=1G
CPUQuota=80%

TimeoutStartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel

[Install]
WantedBy=multi-user.target
